# contrib/buff.service
[Unit]
Description=buff - content relay server
After=network-online.target
Wants=network-online.target

StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=exec
User=buff
Group=buff
ExecStart=/usr/local/bin/buff serve

Environment=BUFF_DATA_DIR=%S/buff
Environment=BUFF_ADDR=:8080
EnvironmentFile=-/etc/buff/buff.env

StateDirectory=buff
StateDirectoryMode=0700

KillSignal=SIGTERM
TimeoutStopSec=20
Restart=on-failure
RestartSec=2

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
CapabilityBoundingSet=
UMask=0077

[Install]
WantedBy=multi-user.target
