# contrib/rc.d/buff
#!/bin/sh
#
# buff — content relay server
#
# Install:
#   pw useradd buff -d /nonexistent -s /usr/sbin/nologin -c "buff content relay"
#   install -m 0755 bin/buff /usr/local/bin/buff
#   install -m 0555 contrib/rc.d/buff /usr/local/etc/rc.d/buff
#
# Then in /etc/rc.conf (or /etc/rc.conf.d/buff):
#   buff_enable="YES"
#   buff_data_dir="/var/db/buff"              # required: the storage root for every clip
#   #buff_addr=":8080"                        # listen address (default :8080)
#   #buff_user="buff"                         # unprivileged service user (default buff)
#   #buff_group="buff"                        # service group (default buff)
#   #buff_env="BUFF_TTL=24h BUFF_FSYNC=on"    # any extra BUFF_* settings, space-separated
#
# PROVIDE: buff
# REQUIRE: NETWORKING DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name="buff"
rcvar="buff_enable"
desc="buff content relay server"

load_rc_config "$name"

: ${buff_enable:="NO"}
: ${buff_addr:=":8080"}
: ${buff_user:="buff"}
: ${buff_group:="buff"}
: ${buff_env:=""}

pidfile="/var/run/${name}.pid"
procname="/usr/local/bin/${name}"

command="/usr/sbin/daemon"
command_args="-S -p ${pidfile} -t ${name} -u ${buff_user} \
  /usr/bin/env BUFF_DATA_DIR=${buff_data_dir} BUFF_ADDR=${buff_addr} ${buff_env} \
  ${procname} serve"

sig_stop="TERM"
start_precmd="buff_prestart"

buff_prestart()
{
	# BUFF_DATA_DIR has no usable default — it is the storage boundary itself — so an unset
	# buff_data_dir fails loudly here, before daemon(8) launches a buff that could only error out.
	if [ -z "${buff_data_dir}" ]; then
		err 1 "buff_data_dir is required: set it in rc.conf to the storage directory"
	fi
	# Create the storage root owned by the unprivileged service user. This runs as root, before
	# daemon(8) drops privileges, because a parent like /var/db is not writable by buff_user, so
	# buff's own mkdir could not create it. install -d is idempotent: an existing dir is left as is.
	if [ ! -d "${buff_data_dir}" ]; then
		install -d -o "${buff_user}" -g "${buff_group}" -m 0700 "${buff_data_dir}" ||
			err 1 "could not create ${buff_data_dir}"
	fi
}

run_rc_command "$1"
